# cookbooks/calendar/templates/opt/kf-calendar/etc/systemd/system/calendar-update-feeds.service.erb
[Unit]
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/kf-calendar/bin/calendar-update-feeds <%= node['user']['kf_calendar']['directory'] %>/yt_api_key.txt
User=kf_calendar
WorkingDirectory=/opt/kf-calendar/share/calendar
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
