# cookbooks/webhook/files/opt/pipitor-infra/etc/systemd/system/webhook.socket
[Unit]
Description=Webhook socket

[Socket]
ListenStream=/run/webhook.sock
SocketUser=nginx
SocketGroup=nginx

[Install]
WantedBy=sockets.target
